# deploy/paper-predict-today.service
[Unit]
Description=Run today's HKJC paper-trade predict
After=network.target

[Service]
Type=oneshot
User=horseracing
Group=horseracing
WorkingDirectory=/home/horseracing/app
Environment=PATH=/home/horseracing/app/venv/bin:/usr/bin
ExecStart=/bin/bash -c '/home/horseracing/app/venv/bin/python scripts/paper_trade.py predict --date $(TZ=Asia/Hong_Kong date +%%Y-%%m-%%d)'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=paper-predict
