# deploy/paper-predict.timer
[Unit]
Description=Daily HKJC paper-trade predict at 17:30 HKT

[Timer]
# 17:30 HKT (09:30 UTC) — racecards are inserted by ~09:00 UTC (17:00 HKT)
OnCalendar=*-*-* 17:30:00 Asia/Hong_Kong
Persistent=true
Unit=paper-predict-today.service

[Install]
WantedBy=timers.target
