# deploy/paper-predict.timer
[Unit]
Description=Daily HKJC paper-trade predict at 11:00 HKT

[Timer]
# 11:00 HKT — racecard is published the evening before and the orchestrator
# re-scrapes it at 09:00 HKT, so predictions land ~2h before the first
# Sunday race (~13:00 HKT) instead of after the card has finished.
OnCalendar=*-*-* 11:00:00 Asia/Hong_Kong
Persistent=true
Unit=paper-predict-today.service

[Install]
WantedBy=timers.target
